# dist/src/main/unix/scripts/postrm
#!/bin/sh


if [ "$1" = "remove" ]; then
  update-rc.d crossdata remove > /dev/null 2>&1 || true
  userdel crossdata
fi